# .local/share/service/picom/run
#!/bin/sh
[ -n "$DISPLAY" ] && exec picom
