# .config/service/gpg-agent/run
#!/bin/sh
exec gpg-agent --server --no-detach
