# .config/service/redshift/run
#!/bin/sh
[ -n "$DISPLAY" ] && exec redshift-gtk
