# .local/share/service/dunst/run
#!/bin/sh
[ -n "$DISPLAY" ] && exec dunst
