# .config/service/pipewire-pulse/run
#!/bin/sh
exec pipewire -c pipewire-pulse.conf
