# .local/share/runit/sv/wireplumber/run
#!/bin/sh

sv check "$HOME/.runit/sv/pipewire" >/dev/null || exit 1

exec wireplumber 2>&1